# modules/discord-adapter/systemd/user/discord-agent-adapter.service
[Unit]
Description=Discord Agent Adapter
Documentation=file:/home/chenma/Documents/My_App_Dev/discord_agent_adapter/README.md
After=network-online.target agent-host-web.service
Wants=network-online.target agent-host-web.service

[Service]
Type=simple
WorkingDirectory=/home/chenma/Documents/My_App_Dev/discord_agent_adapter
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=%h/.config/agent-host/secrets.env
ExecStartPre=/home/chenma/Documents/My_App_Dev/discord_agent_adapter/.venv/bin/python /home/chenma/Documents/My_App_Dev/discord_agent_adapter/bot.py --config /home/chenma/Documents/My_App_Dev/discord_agent_adapter/config.json --check-config
ExecStart=/home/chenma/Documents/My_App_Dev/discord_agent_adapter/.venv/bin/python /home/chenma/Documents/My_App_Dev/discord_agent_adapter/bot.py --config /home/chenma/Documents/My_App_Dev/discord_agent_adapter/config.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
